# systemd/xsmom-optimizer-full-cycle.service
[Unit]
Description=XSMOM Full-Cycle Optimizer (WFO + BO + MC)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/xsmom-bot
EnvironmentFile=/opt/xsmom-bot/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/xsmom-bot
Environment=PATH=/opt/xsmom-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/bin/bash -lc '/opt/xsmom-bot/bin/run-optimizer-full-cycle.sh /opt/xsmom-bot/config/config.yaml'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xsmom-optimizer-full-cycle

# Full-cycle optimizer can take 3-6 hours depending on:
# - Number of BO trials (default: 100)
# - Number of MC runs (default: 1000)
# - Number of WFO segments
# - Data size and API response times
# Set timeout to 6 hours (21600 seconds) to allow completion
TimeoutStartSec=21600

[Install]
WantedBy=multi-user.target

