# systemd/xsmom-optimizer.timer
[Unit]
Description=Nightly XSMOM optimizer run

[Timer]
# Runs every day at 03:30 local time with up to 20 minutes of jitter
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=20m
Persistent=true
Unit=xsmom-optimizer.service

[Install]
WantedBy=timers.target
