# systemd/xsmom-opt.service
[Unit]
Description=XSMOM Optimizer (timeframe + regime) — updates config and restarts bot if improved
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/xsmom-bot
EnvironmentFile=/opt/xsmom-bot/.env
ExecStart=/opt/xsmom-bot/venv/bin/python -m src.auto_opt \
  --config /opt/xsmom-bot/config/config.yaml \
  --service xsmom-bot.service \
  --restart
Nice=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
